# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:openvscode-server]
command=/usr/bin/openvscode-server --host 127.0.0.1 --port 8000 --without-connection-token --default-folder /data/workspace --user-data-dir /data/.openvscode --extensions-dir /data/.openvscode/extensions
priority=100
autostart=true
autorestart=true
stdout_logfile=/var/log/openvscode.log
stderr_logfile=/var/log/openvscode.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=VSCODE_SETTINGS="/app/config/vscode-settings.json"

[program:kernel-server]
command=tsx /app/kernel/server.ts
directory=/app
priority=200
autostart=true
autorestart=true
stdout_logfile=/var/log/kernel.log
stderr_logfile=/var/log/kernel.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
